# dot_config/systemd/user/rclone-Documents.service
[Unit]
Description=rclone File Synchronization

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync "home:/Documents/" "onedrive:/Documents/" --create-empty-src-dirs --compare size,modtime,checksum --resilient --recover -MvP --fix-case
