# dot_config/systemd/user/rclone-Pictures.timer
[Unit]
Description=rclone File Synchronization

[Timer]
OnBootSec=10m
OnUnitActiveSec=120m

[Install]
WantedBy=timers.target
