# dot_config/systemd/user/rclone-Videos.timer
[Unit]
Description=rclone File Synchronization

[Timer]
OnBootSec=10m
OnUnitActiveSec=69m

[Install]
WantedBy=timers.target
