# private_dot_config/systemd/user/unison@.service
[Unit]
Description=Unison File Synchronization

[Service]
Type=oneshot
ExecStart=/usr/bin/unison %i
